# sv/influxdb/run
#!/bin/sh

#sv start socklog-unix || exit 1

exec 2>&1
exec influxdb -config /opt/influxdb/shared/config.toml 
